# systemd/obd2_startup.service
[Unit]
Description=OBD2 Startup Service
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/obd2/obd2-repo
ExecStart=/opt/obd2/venv/bin/python /opt/obd2/obd2-repo/scripts/system_check_fix.py --fix
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
