# systemd/obd2_hub.service
[Unit]
Description=OBD2 Hub Service (Disabled for ESP32 setup)
After=network.target obd2_startup.service

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
